# Lesson17/provisioning/files/borg-backup.timer
[Unit]
Description=Borg Backup
Requires=borg-backup.service
[Timer]
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target

